# system_files/etc/systemd/system/usb-monitor.service
[Unit]
Description=USB monitor notifications

[Service]
ExecStart=/usr/libexec/usb-monitor.sh
Restart=always

[Install]
WantedBy=default.target